# docker/php-config/supervisord.conf
[supervisord]
nodaemon=true                                 ; start in foreground if true; default false
logfile=/etc/supervisor/logs/supervisord.log  ; main log file; default $CWD/supervisord.log
logfile_maxbytes=5MB                          ; max main logfile bytes; default 50MB
logfile_backups=0                             ; # of main logfile backups; 0 means none, default 10
user=root

[supervisorctl]

[inet_http_server]
port=127.0.0.1:9001

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[group:worker]
programs=queue-worker,schedule-worker

[program:php-fpm]
command=php-fpm
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
autostart=true
autorestart=true
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:queue-worker]
command=php artisan queue:work
directory=/srv
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
autostart=true
autorestart=true
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:schedule-worker]
command=php artisan schedule:work
directory=/srv
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
autostart=true
autorestart=true
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
